# config/rc.local
#! /bin/bash

root_device=$(findmnt / --noheadings --output SOURCE)

if [ -e /etc/configboot ]; then

  hostnamectl set-hostname rock4se
  chmod u+s /usr/bin/sudo
  dpkg-reconfigure lightdm
  dpkg-reconfigure sddm
  dpkg-reconfigure gdm3
  rm /etc/configboot
  touch /etc/firstboot
  poweroff
fi

if [ -e /etc/firstboot ]; then
  resize2fs $root_device
  rm /etc/firstboot
fi
exit 0