# config/rc.local
#! /bin/bash

if [ -e /etc/configboot ]; then
  hostnamectl set-hostname rock4se
  chmod u+s /usr/bin/sudo
  dpkg-reconfigure lightdm
  dpkg-reconfigure sddm
  dpkg-reconfigure gdm3
  touch /etc/default/FIRSTBOOT
  resizeroot
  rm /etc/configboot
  reboot
fi

exit 0
