# config/rc.local
#! /bin/bash

if [ -e /etc/configboot ]; then
  hostnamectl set-hostname rock4se
  chmod u+s /usr/bin/sudo
  dpkg-reconfigure lightdm
  dpkg-reconfigure sddm
  dpkg-reconfigure gdm3
  rm /etc/configboot
  touch /etc/default/FIRSTBOOT
  poweroff
fi

if [ -e /etc/default/FIRSTBOOT ]; then
  /etc/init.d/resize2fs start
  rm /etc/firstboot
fi
exit 0
