# config/rc.local
#! /bin/bash

if [ -f /etc/default/FIRSTBOOT ]; then
  /etc/init.d/resize2fs start
fi

if [ -e /etc/configboot ]; then
  hostnamectl set-hostname rock4se
  chmod u+s /usr/bin/sudo
  if dpkg -l | grep -q "lightdm"; then
  	dpkg-reconfigure lightdm
  fi
  if dpkg -l | grep -q "sddm"; then
  	dpkg-reconfigure sddm
  fi
  if dpkg -l | grep -q "gdm3"; then
  	dpkg-reconfigure gdm3
  fi
  rm /etc/configboot
  touch /etc/default/FIRSTBOOT
  poweroff
fi

exit 0
